# debian/somatic-utils.postrm
update-rc.d somatic remove
#DEBHELPER#

